# docs/systemd/attendx-backend.service
[Unit]
Description=AttendX Backend (FastAPI)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/attendx
ExecStart=/usr/bin/docker compose -f docker/docker-compose.prod.yml up backend
ExecStop=/usr/bin/docker compose -f docker/docker-compose.prod.yml stop backend
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=attendx-backend

[Install]
WantedBy=multi-user.target
