# docs/systemd/attendx-bot.service
[Unit]
Description=AttendX Telegram Bot
After=network.target docker.service attendx-backend.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/attendx
ExecStart=/usr/bin/docker compose -f docker/docker-compose.prod.yml up bot
ExecStop=/usr/bin/docker compose -f docker/docker-compose.prod.yml stop bot
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=attendx-bot

[Install]
WantedBy=multi-user.target
